# rhq/image/start.ini
[supervisord]
nodaemon=true

[program:postgresql]
priority=1
command=runuser -l postgres -c "pg_ctl -D /var/lib/pgsql/data/ -l logfile start"

[program:rhq]
priority=100
command=/rhq-server-4.11.0/bin/rhqctl start

[program:sshd]
command=/usr/sbin/sshd -D

