# infinispan-spark/supervisord.conf
[supervisord]
nodaemon=true

[program:spark_master]
command=/usr/local/spark/start-spark.sh

[program:infinispan]
command=/usr/local/infinispan-server/bin/start-infinispan.sh



